# system/fpp-monitor-agent.service
[Unit]
Description=FPP Monitor Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fpp
Group=fpp
ExecStart=/opt/fpp-monitor-agent/fpp-monitor-agent
Restart=always
RestartSec=5
Environment=FPP_MONITOR_AGENT_CONFIG=/home/fpp/media/config/fpp-monitor-agent.json
# Logs go to journald by default

[Install]
WantedBy=multi-user.target
